# django_deploy/polls.ini
[program:polls]
command=/root/Desktop/python-envl/django1.8/bin/python /project/baidu/tornado_server.py
directory=/project/baidu